# word_banker_api/run1.conf
[program:api1]directory=/home/word_banker_apicommand=python3 -m uvicorn main:app --reload autostart=trueautorestart=trueenvironment=HOME="/home",USER="root"stderr_logfile=/home/word_banker_api/logfile_err.logstdout_logfile=/home/word_banker_api/logfile.log